# podman/immich/quadlet/immich-redis.container
[Container]
ContainerName=redis
HealthCmd=redis-cli ping || exit 1
Image=docker.io/valkey/valkey:8-bookworm
Pod=immich.pod

[Service]
Restart=always
